# packages/luci-app-h5000m-netmode/root/etc/uci-defaults/90-h5000m-netmode
#!/bin/sh

/usr/sbin/h5000m-netmode >/dev/null 2>&1 || true

exit 0
