# files/etc/uci-defaults/99-h5000m-defaults
#!/bin/sh

logger -t h5000m-defaults "Set default timezone and locale defaults"
uci -q set system.@system[0].zonename='Asia/Shanghai'
uci -q set system.@system[0].timezone='CST-8'
uci -q get system.ntp >/dev/null 2>&1 || uci -q set system.ntp='timeserver'
uci -q set system.ntp.enabled='1'
uci -q delete system.ntp.server
uci -q add_list system.ntp.server='ntp.aliyun.com'
uci -q add_list system.ntp.server='cn.pool.ntp.org'
uci -q add_list system.ntp.server='time1.cloud.tencent.com'
uci -q add_list system.ntp.server='time.ustc.edu.cn'
uci -q commit system

if [ -f /etc/config/luci ]; then
	uci -q get luci.main >/dev/null 2>&1 || uci -q set luci.main='core'
	uci -q set luci.main.lang='zh_cn'
	uci -q commit luci
fi

if [ -f /etc/config/uhttpd ]; then
	logger -t h5000m-defaults "Redirect LuCI management traffic to HTTPS"
	uci -q set uhttpd.main.redirect_https='1'
	uci -q commit uhttpd
fi

if [ -f /etc/config/dropbear ]; then
	logger -t h5000m-defaults "Disable SSH password authentication by default"
	uci -q set dropbear.@dropbear[0].PasswordAuth='off'
	uci -q set dropbear.@dropbear[0].RootPasswordAuth='off'
	uci -q commit dropbear
fi

if [ -x /etc/init.d/ttyd ]; then
	logger -t h5000m-defaults "Disable Web terminal by default"
	/etc/init.d/ttyd disable >/dev/null 2>&1 || true
fi

if [ -e /sys/class/leds/mdio-bus:0f:red:wan ]; then
	logger -t h5000m-defaults "Bind the red WAN LED to eth1 link and traffic"
	uci -q set system.led_wan='led'
	uci -q set system.led_wan.name='WAN'
	uci -q set system.led_wan.sysfs='mdio-bus:0f:red:wan'
	uci -q set system.led_wan.trigger='netdev'
	uci -q set system.led_wan.dev='eth1'
	uci -q delete system.led_wan.mode
	uci -q add_list system.led_wan.mode='link'
	uci -q add_list system.led_wan.mode='tx'
	uci -q add_list system.led_wan.mode='rx'
	uci -q commit system
fi

logger -t h5000m-defaults "Set default LAN address to 192.168.10.1"
uci -q set network.lan.proto='static'
uci -q set network.lan.ipaddr='192.168.10.1'
uci -q set network.lan.netmask='255.255.255.0'
uci -q set network.wan.norelease='1'

logger -t h5000m-defaults "Write default H5000M QModem slot"
if [ -f /etc/config/qmodem ]; then
	# The upstream OpenWrt device tree enumerates the built-in MT5700M on
	# usb1/1-1, while older H5000M builds exposed the same hardware as 2-1.
	# Prefer the currently enumerated modem and fall back to the upstream slot.
	qmodem_slot='1-1'
	for slot in 1-1 2-1; do
		usb_path="/sys/bus/usb/devices/${slot}"
		[ -d "${usb_path}" ] || continue
		usb_vendor="$(cat "${usb_path}/idVendor" 2>/dev/null || true)"
		usb_product="$(cat "${usb_path}/product" 2>/dev/null || true)"
		case "${usb_vendor}:${usb_product}" in
			3466:*|*:[Mm][Tt]5700*)
				qmodem_slot="${slot}"
				break
				;;
		esac
	done
	uci -q set qmodem.usb0='modem-slot'
	uci -q set qmodem.usb0.slot="${qmodem_slot}"
	uci -q set qmodem.usb0.type='usb'
	uci -q set qmodem.usb0.alias='MT5700M'
	uci -q set qmodem.main.enable_pcie_scan='0'
	uci -q set qmodem.main.enable_usb_scan='0'
	uci -q set qmodem.main.try_preset_pcie='0'
	uci -q set qmodem.main.try_preset_usb='1'
	uci -q set qmodem.main.start_delay='0'
	uci -q commit qmodem
fi

logger -t h5000m-defaults "Apply WAN and 5G default route policy"
if [ -x /usr/sbin/h5000m-netmode ]; then
	/usr/sbin/h5000m-netmode >/dev/null 2>&1 || true
fi

logger -t h5000m-defaults "Remove build-only apk feed entries"
if [ -f /etc/apk/repositories.d/distfeeds.list ]; then
	sed -i '/\/qmodem\//d;/ qmodem[ _-]/d;/ qmodem$/d;/\/small_package\//d;/ small_package[ _-]/d;/ small_package$/d;/\/video\//d' /etc/apk/repositories.d/distfeeds.list
fi

exit 0
