# files/etc/hotplug.d/iface/99-h5000m-wan-priority
#!/bin/sh

[ "${ACTION}" = "ifup" ] || [ "${ACTION}" = "ifupdate" ] || exit 0

case "${INTERFACE}" in
	wan|wan6|USB|USBv6)
		;;
	*)
		exit 0
		;;
esac

logger -t h5000m-wan-priority "同步 WAN 与 5G 备用出口优先级"

if [ -x /usr/sbin/h5000m-netmode ]; then
	/usr/sbin/h5000m-netmode >/dev/null 2>&1 || true
fi
